# overlays/odin-debugging/etc/rc.local
#!/bin/bash
echo "*** rc.local ***"
set -x
ip addr
networkctl
networkctl status